# deploy/mibee-steward.service
[Unit]
Description=MiBee Steward - Device Management & Monitoring
After=network.target
Wants=network-online.target

[Service]
Type=simple
# Use the NUMERIC uid/gid, not the username. On Arch/systemd ≥256 with
# `passwd: files systemd` in nsswitch.conf, systemd resolves User= via its own
# nss-systemd module, which intermittently fails with "Unknown user" (status
# 217/USER) for perfectly valid local users right after a daemon-reload. The
# numeric form bypasses NSS entirely. 1000 = the mibee service account.
User=1000
Group=1000
WorkingDirectory=/opt/mibee-steward
ExecStart=/opt/mibee-steward/mibee-steward -config configs/config.yaml
# always (not on-failure): a monitoring daemon should come back after ANY exit,
# including a clean one (e.g. an unexpected SIGTERM from outside). `systemctl
# stop` is still honored — systemd won't restart after an admin stop. Without
# this, a graceful exit leaves the service down for good.
Restart=always
# 10s (not the default 5s) gives the kernel time to release the previous
# process's TCP socket (TIME_WAIT) before the restart attempt — without it the
# new process hits "bind: address already in use" and the cycle repeats. The
# in-process bind retry (listenAndServeWithRetry in main.go) is the primary
# guard; this is defense-in-depth.
RestartSec=10
# No StartLimitBurst: a monitoring daemon must NOT be rate-limited into
# permanent "failed" state by a transient bind storm — that leaves it down
# until manual intervention. Restart=always with no burst cap means systemd
# keeps trying indefinitely (backed off by RestartSec).
LimitNOFILE=65536
StandardOutput=journal
StandardError=journal
LogRateLimitIntervalSec=30s
LogRateLimitBurst=1000
MemoryMax=1G

# Optional environment overrides via /etc/default/mibee-steward
EnvironmentFile=-/etc/default/mibee-steward

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/mibee-steward/data
PrivateTmp=true
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
